# templates/ddns.timer.j2
[Unit]
Description=Update DDNS (timer)

[Timer]
OnUnitActiveSec={{ ddns__update_interval_minutes }}m
RandomizedDelaySec=1m

[Install]
WantedBy=timers.target
